# configs/supervisor.d/cesi.ini
[program:cesi]
command=/usr/bin/python2 /app/cesi/web.py
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cesi.log
